# conexus/daemon/com.nexus.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the nexus storage SERVICE (RDR-174 P2.1, nexus-y2yj6).
  Installed by ``nx daemon service install --autostart``. The ``__NX_BIN__``
  placeholder is substituted at install time with the resolved absolute path
  to the ``nx`` executable.

  ``nx daemon service start --foreground`` runs the storage supervisor
  (run_storage_supervisor) in the foreground — it owns the native engine-
  service binary + the local Postgres cluster and blocks until SIGTERM, so
  launchd's supervision is the canonical lifecycle owner.

  No PG readiness wrapper (RDR-174 P2.2, nexus-exfns — VERIFIED NO-OP): the
  supervisor starts its OWN nx-owned Postgres cluster as step 1 of startup
  (storage_service_daemon._start_locked -> _ensure_pg_running ->
  pg_provision._start_cluster), with boot-safe binary discovery from the config
  dir. launchd having no After= equivalent does not matter — there is no
  external Postgres to wait for; a wrapper would poll a PG nothing else starts.
  RunAtLoad launches the supervisor directly and it brings up its own PG.

  StandardOutPath / StandardErrorPath are the launchd-captured supervisor
  stream (interpreter/boot crash diagnostics). The bounded primary log is
  ``~/.config/nexus/logs/`` written by the supervisor itself.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NX_BIN__</string>
        <string>daemon</string>
        <string>service</string>
        <string>start</string>
        <string>--foreground</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!--
      restart-on-FAILURE only (nexus-6bmph, RDR-183 defect-3; parity with the
      systemd unit's Restart=on-failure): SuccessfulExit=<false/> means launchd
      respawns the supervisor only when it exits NON-zero. A SUCCESSFUL exit
      settles the unit, which is load-bearing twice over:
        - healthy coexistence: start() finds a live session-started supervisor
          already owning the lease and exits 0 (exit_if_process_unowned) — the
          old KeepAlive=<true/> respawned that exit-0 every ThrottleInterval
          FOREVER (the GH #1405 defect-3 steady-state churn);
        - a graceful ``nx daemon service stop`` (SIGTERM -> clean exit 0) now
          STAYS stopped instead of being resurrected within 30s.
      Failure exits (crash; exit 2 pre-provisioning before ``nx init``; exit 3
      watchdog service-unrecoverable) still respawn indefinitely — the
      never-give-up posture is preserved exactly where it belongs. Boot
      coverage is RunAtLoad; if a session-owned supervisor later dies, the
      client-side heal-on-next-use respawn path recovers it (RDR-175's OS
      watchdog supervises processes THIS unit started, not foreign leases).
      ExitTimeOut bounds a stuck graceful shutdown.

      ThrottleInterval (RDR-174 P2.1 review SIG-1): bounds the FAILURE respawn
      loop (e.g. crash-looping before ``nx init`` provisions pg_credentials).
      30s mirrors the Linux RestartSec spirit at a safer interval.
    -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>ExitTimeOut</key>
    <integer>20</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/nexus-service.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/nexus-service.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH_ENV__</string>
    </dict>
</dict>
</plist>
