# conexus/daemon/nexus-service.service
## systemd user unit for the nexus storage SERVICE (RDR-174 P2.1, nexus-y2yj6).
## Installed by ``nx daemon service install --autostart``. The ``__NX_BIN__``
## placeholder is substituted at install time with the resolved absolute path
## to the ``nx`` executable.
##
## ``nx daemon service start --foreground`` runs the storage supervisor
## (run_storage_supervisor) in the foreground — it owns the native engine-
## service binary + the local Postgres cluster and blocks until SIGTERM, so
## systemd's supervision is the canonical lifecycle owner.
##
## RDR-174 P2.2 will add the Postgres boot-ordering delta here
## (After=postgresql.service + Wants=postgresql.service). This template is
## intentionally the plain unit; that ordering seam lands in P2.2.

[Unit]
Description=Nexus storage service (engine-service binary + local Postgres; serves T2 + T3)
After=network.target
PartOf=default.target

[Service]
Type=simple
ExecStart=__NX_BIN__ daemon service start --foreground
Restart=on-failure
RestartSec=5s
## SIGTERM-induced exits surface as code 143. Mark them as success so a
## graceful ``nx daemon service stop`` (which sends SIGTERM) does not
## immediately respawn through Restart=on-failure.
SuccessExitStatus=143
Environment=PATH=__PATH_ENV__
StandardOutput=append:__LOG_DIR__/nexus-service.log
StandardError=append:__LOG_DIR__/nexus-service.err

[Install]
WantedBy=default.target
