# nx/daemon/com.nexus.t2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the nexus T2 daemon (RDR-112, nexus-6w0c).
  Installed by ``nx daemon t2 install --autostart``. The ``__NX_BIN__``
  placeholder is substituted at install time with the resolved absolute
  path to the ``nx`` executable.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.t2</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NX_BIN__</string>
        <string>daemon</string>
        <string>t2</string>
        <string>start</string>
        <string>--foreground</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/nexus-t2.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/nexus-t2.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH_ENV__</string>
    </dict>
</dict>
</plist>
