# nx/daemon/nexus-t2.service
## systemd user unit for the nexus T2 daemon (RDR-112, nexus-6w0c).
## Installed by ``nx daemon t2 install --autostart``. The ``__NX_BIN__``
## placeholder is substituted at install time with the resolved absolute
## path to the ``nx`` executable.

[Unit]
Description=Nexus T2 daemon (memory + tuplespace)
After=network.target
PartOf=default.target

[Service]
Type=simple
ExecStart=__NX_BIN__ daemon t2 start --foreground
Restart=on-failure
RestartSec=5s
Environment=PATH=__PATH_ENV__
StandardOutput=append:__LOG_DIR__/nexus-t2.log
StandardError=append:__LOG_DIR__/nexus-t2.err

[Install]
WantedBy=default.target
